# monitoring/lanxi-monitor/systemd/lanxi-monitor.service
[Unit]
Description=LANXI Monitor Service
After=network.target

[Service]
ExecStart=/home/ubuntu/lanxi/bin/lanxi
WorkingDirectory=/home/ubuntu/lanxi
Restart=always
User=pi
Group=pi
Environment="PORT=8080"
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
